# 99-external-display.rules
SUBSYSTEM=="drm", ACTION=="add", RUN+="/home/development/.resize-windows.sh" OPTIONS="last_rule"
